# ubuntuserver/systemd/user/hay-say-runtime-manager.service
[Unit]
Description=Hay Say model runtime manager
PartOf=hay-say.target
Requires=hay-say-redis.service
After=hay-say-redis.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/hay-say/environment
ExecStart=/usr/bin/env %h/.local/lib/hay-say/bin/run-service runtime-manager
ExecStop=-/usr/bin/curl --fail --silent --show-error --max-time 25 --request POST http://127.0.0.1:6588/runtimes/stop-all
Restart=on-failure
RestartSec=3
TimeoutStopSec=30
KillMode=control-group
TasksMax=4096
LimitNOFILE=65536
UMask=0077
